# private_dot_config/systemd/user/clean-go-caches.timer
[Unit]
Description=Clean Go cache folders monthly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=default.target
