# private_dot_config/systemd/user/update-dotfiles.timer
[Unit]
Description=Update dotfiles on boot

[Timer]
OnBootSec=10m

[Install]
WantedBy=timers.target
